# deploy/yardmaster.service
[Unit]
Description=YardMaster Django App
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/yardmaster
EnvironmentFile=/var/www/yardmaster/.env
ExecStart=/var/www/yardmaster/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8000 \
    yardmaster.wsgi:application
Restart=always

[Install]
WantedBy=multi-user.target
